# infrastructure/la-egress-lockdown/systemd/n8n.service.d/proxy.conf
# BHN — n8n outbound HTTP via Hillsboro tinyproxy.
# Drop-in at /etc/systemd/system/n8n.service.d/proxy.conf
# After dropping the file:  systemctl daemon-reload && systemctl restart n8n
#
# n8n's HTTP-request nodes use Node.js's built-in HTTP/HTTPS modules, which
# honor HTTPS_PROXY/HTTP_PROXY env vars. Anthropic SDK in n8n also honors them.
#
# NO_PROXY keeps intra-mesh + local traffic off the proxy path — important
# for n8n→PG (localhost), n8n→FRA (10.9.0.0/24), n8n→Hillsboro mgmt
# (10.8.0.0/24 covers it), and the actual proxy host (10.8.0.6) so we don't
# try to proxy the proxy.

[Service]
Environment=HTTP_PROXY=http://10.8.0.6:8888
Environment=HTTPS_PROXY=http://10.8.0.6:8888
Environment=NO_PROXY=localhost,127.0.0.1,10.8.0.0/24,10.9.0.0/24
