# infrastructure/la-egress-lockdown/systemd/docker.service.d/proxy.conf
# BHN — Docker daemon outbound HTTP via Hillsboro tinyproxy.
# Drop-in at /etc/systemd/system/docker.service.d/proxy.conf on LA.
#
# Affects: docker pull, docker push, docker registry lookups, anything the
# DAEMON itself talks to over HTTP/HTTPS. Does NOT auto-pass into containers
# at runtime — containers inherit their own env via `docker run -e` or
# compose env: blocks.
#
# Activation: takes effect after `systemctl daemon-reload && systemctl
# restart docker`. The restart briefly stops all running containers
# (~5-15 sec). Coordinate timing with operator — LA currently runs
# n8n, bhn-horizon-redis, bhn-wallos. Restart only at a safe window.
#
# After dropping the file and reloading systemd, verify:
#   systemctl show docker | grep -i environment
#   docker info | grep -i proxy

[Service]
Environment="HTTP_PROXY=http://10.8.0.6:8888"
Environment="HTTPS_PROXY=http://10.8.0.6:8888"
Environment="NO_PROXY=localhost,127.0.0.1,10.8.0.0/24,10.9.0.0/24,docker0,*.docker"
