# infrastructure/la-egress-lockdown/systemd/grafana-server.service.d/proxy.conf
# BHN — Grafana outbound HTTP via Hillsboro tinyproxy.
# Drop-in at /etc/systemd/system/grafana-server.service.d/proxy.conf
# After dropping:  systemctl daemon-reload && systemctl restart grafana-server
#
# Grafana's alert webhook deliveries (currently to n8n at <BHN_WG_LA_IP>:5678 —
# intra-mesh, kept off proxy via NO_PROXY) and plugin-update calls
# (out to grafana.com — these should go via proxy).

[Service]
Environment=HTTP_PROXY=http://<BHN_WG_HIL_IP>:8888
Environment=HTTPS_PROXY=http://<BHN_WG_HIL_IP>:8888
Environment=NO_PROXY=localhost,127.0.0.1,10.8.0.0/24,10.9.0.0/24
