# scripts/weather-collectors/systemd-units/bhn-weather-collector-helsinki.timer
# BHN WeatherBHN — Helsinki collector cadence. Matches LA's original
# bhn-weather-collector.timer (every 30 minutes).

[Unit]
Description=BHN WeatherBHN collector — Helsinki (every 30 minutes)

[Timer]
OnCalendar=*-*-* *:00,30:00
Unit=bhn-weather-collector-helsinki.service
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
