# scripts/trading/systemd-units/bhn-weather-collector.service
# BHN trading — Strat 9 (BHN-WEATHER-ALPHA) data collector — LA (hub).
#
# 2026-07-01: Open-Meteo, ASOS, NWS, NOAA/NBM, NOMADS, USDA-crops sources
# moved to dedicated collector nodes (Helsinki: open_meteo*; Hillsboro:
# asos/nws*/nbm/nomads/usda_crops — see
# scripts/weather-collectors/systemd-units/). LA now only runs the
# auth-sensitive Kalshi sources, which stay on LA by design (Kalshi API
# keys never leave the hub).
#
# Phase 1 — data collection only. No betting, no strategy logic.
# Scheduled by bhn-weather-collector.timer (every 30 minutes).
#
[Unit]
Description=BHN Strat 9 weather data collector — LA (kalshi_markets + kalshi_portfolio only)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/bhn/trading

# Trading framework env (PG_HOST/PG_PORT/PG_DB/PG_USER/PG_PASSWORD,
# plus KALSHI_API_KEY for the market/portfolio sources)
EnvironmentFile=/etc/bhn-trading/env
EnvironmentFile=-/etc/bhn-trading/strat9.env

ExecStart=/bin/bash -c 'for s in kalshi_markets kalshi_portfolio; do /usr/bin/python3 /opt/bhn/trading/weather_data_collector.py --source "$s" || exit 1; done'

StandardOutput=append:/var/log/bhn-trading/weather-collector.log
StandardError=append:/var/log/bhn-trading/weather-collector.log

TimeoutStartSec=300
TimeoutStopSec=30

[Install]
# Started by the timer; not directly enabled.
