# scripts/trading/systemd-units/bhn-weather-collector.timer
# BHN trading — Strat 9 weather collector cadence.
# Every 30 minutes.

[Unit]
Description=BHN Strat 9 weather collector (every 30 minutes)

[Timer]
OnCalendar=*-*-* *:00,30:00
Unit=bhn-weather-collector.service
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
