# scripts/trading/systemd-units/bhn-kalshi-portfolio.service
# BHN trading — Strat 9 Kalshi portfolio snapshot collector.
# Calls fetch_kalshi_portfolio(): GET /portfolio/positions + /portfolio/fills
# and writes to kalshi_positions + kalshi_fills tables.
#
# Scheduled by bhn-kalshi-portfolio.timer (every 15 minutes).
#
# Deploy:
#   scp scripts/trading/systemd-units/bhn-kalshi-portfolio.service root@10.8.0.1:/etc/systemd/system/
#   scp scripts/trading/systemd-units/bhn-kalshi-portfolio.timer   root@10.8.0.1:/etc/systemd/system/
#   ssh root@10.8.0.1 "systemctl daemon-reload && systemctl enable --now bhn-kalshi-portfolio.timer"

[Unit]
Description=BHN Strat 9 Kalshi portfolio snapshot (positions + fills)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/bhn/trading

EnvironmentFile=/etc/bhn-trading/env
EnvironmentFile=-/etc/bhn-trading/strat9.env

ExecStart=/usr/bin/python3 /opt/bhn/trading/weather_data_collector.py --source kalshi_portfolio

StandardOutput=append:/var/log/bhn-trading/kalshi-portfolio.log
StandardError=append:/var/log/bhn-trading/kalshi-portfolio.log

TimeoutStartSec=60
TimeoutStopSec=15

[Install]
# Started by the timer; not directly enabled.
