# scripts/trading/systemd-units/bhn-weather-position-monitor.timer
[Unit]
Description=BHN WeatherBHN stop-loss position monitor (every 60 seconds)

[Timer]
OnCalendar=*-*-* *:*:00
Unit=bhn-weather-position-monitor.service
# Persistent=false intentional: don't replay missed fires on startup —
# stop-loss checks are only meaningful against live market data.
Persistent=false

[Install]
WantedBy=timers.target
