# infrastructure/grafana/systemd/grafana-server.service.d/nj-env.conf
# /etc/systemd/system/grafana-server.service.d/nj-env.conf
#
# Loads /etc/bhn-trading/env (mode 0600) which carries:
#   GF_DATABASE_GRAFANA_READER_PASSWORD=<grafana_reader PG password>
#
# Per the BHN data pipeline standard, datasource passwords are NEVER
# committed to the repo. They live in /etc/bhn-trading/env on each host
# that needs them. The systemd unit injects them into Grafana's process
# environment so the provisioned datasource YAML can reference them via
# $GF_DATABASE_GRAFANA_READER_PASSWORD.

[Service]
EnvironmentFile=/etc/bhn-trading/env
